# com.garrettsmith.daily_worker.py.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
Template LaunchAgent for the Daily Worker auto-form script.

Usage (replace [YOUR_USERNAME]):

    cp launchd/com.example.dailyworker.plist \
       ~/Library/LaunchAgents/com.[YOUR_USERNAME].daily_worker.py.plist

    # Edit the copied file:
    #   - Set Label
    #   - Set ProgramArguments python path and script path
    #   - Set Hour/Minute as desired

Then:

    launchctl bootout gui/$(id -u) ~/Library/LaunchAgents/com.[YOUR_USERNAME].daily_worker.py.plist 2>/dev/null || true
    launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.[YOUR_USERNAME].daily_worker.py.plist

-->
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.YOUR_USERNAME.daily_worker.py</string>

    <key>ProgramArguments</key>
    <array>
      <!-- Adjust Python path and repo path for your system -->
      <string>/usr/local/bin/python3.11</string>
      <string>/Users/YOUR_USERNAME/Projects/dailyworker/run_daily_worker.py</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key>
      <integer>6</integer>
      <key>Minute</key>
      <integer>30</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/YOUR_USERNAME/Library/Logs/run_daily_worker.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USERNAME/Library/Logs/run_daily_worker.err</string>

    <!-- For testing: run immediately on load
         Remove or set to false in final use -->
    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>
  </dict>
</plist>
